# farmboy-schedule-creator.timer
[Unit]
Description=Farm Boy Schedule Creator

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
